# nfs-bind/mnt-nas_media.automount
[Unit]
Description=NAS NFS Media Automount

[Automount]
Where=/mnt/nas_media
TimeoutIdleSec=6000

[Install]
WantedBy=multi-user.target
